# debian/arduino-router/DEBIAN/postinst
#!/bin/sh

systemctl enable arduino-router
systemctl enable arduino-router-serial

