# debian/arduino-router/DEBIAN/postinst
#!/bin/sh

if [ -d /run/systemd/system ]; then
  systemctl daemon-reload
fi

systemctl enable arduino-router
systemctl enable arduino-router-serial

