# debian/arduino-router/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-router
systemctl disable arduino-router-serial

