# debian/router/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-router

if [ -d /run/systemd/system ]; then
  systemctl stop arduino-router
fi
